# rootfs/etc/s6-overlay/s6-rc.d/nut-server/run
#!/command/execlineb -P

with-contenv
foreground { /usr/sbin/upsdrvctl -u root start }
foreground { s6-setuidgid nut /usr/sbin/upsd -D }
